# automation/com.s4ai.mayberry-google-reviews.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.s4ai.mayberry-google-reviews</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>-lc</string>
    <string>cd '/Users/rajatsingh/Documents/Mayberry Pressure Washing llc' &amp;&amp; /usr/bin/env node scripts/sync-google-reviews.mjs --update-site --reply-unanswered --maybe-post &gt;&gt; logs/review-automation.log 2&gt;&amp;1</string>
  </array>
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>7</integer>
    <key>Minute</key>
    <integer>15</integer>
  </dict>
  <key>StandardOutPath</key>
  <string>/Users/rajatsingh/Documents/Mayberry Pressure Washing llc/logs/review-automation-launchd.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/rajatsingh/Documents/Mayberry Pressure Washing llc/logs/review-automation-launchd.err</string>
  <key>RunAtLoad</key>
  <true/>
</dict>
</plist>
